# launchd/com.seanwang.us-premarket-brief.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.seanwang.us-premarket-brief</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/sean.wang/us-premarket-brief/scripts/run-brief.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key><integer>1</integer>
      <key>Hour</key><integer>20</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>2</integer>
      <key>Hour</key><integer>20</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>3</integer>
      <key>Hour</key><integer>20</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>4</integer>
      <key>Hour</key><integer>20</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>5</integer>
      <key>Hour</key><integer>20</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/sean.wang/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/sean.wang</string>
    <key>LANG</key>
    <string>en_US.UTF-8</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/sean.wang/Library/Logs/us-premarket-brief.launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/sean.wang/Library/Logs/us-premarket-brief.launchd.log</string>

  <key>RunAtLoad</key>
  <false/>

  <key>WorkingDirectory</key>
  <string>/Users/sean.wang/us-premarket-brief</string>
</dict>
</plist>
